# ops/systemd/sanctum-mcp-health.timer
[Unit]
Description=Sanctum MCP Health Check Timer (every 30s)

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
Unit=sanctum-mcp-health.service

[Install]
WantedBy=timers.target
