# ops/systemd/sanctum-worker.service
[Unit]
Description=Sanctum Signal Worker (Batch Processor)
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/sanctum-core
EnvironmentFile=/var/www/sanctum-core/.env
ExecStart=/var/www/sanctum-core/venv/bin/python3 -m app.worker

[Install]
WantedBy=multi-user.target
