# ops/systemd/sanctum-mcp.service
[Unit]
Description=Sanctum MCP Server (FastMCP)
After=network.target sanctum-api.service
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
User=preginald
Group=preginald
WorkingDirectory=/home/preginald/DigitalSanctum/sanctum-mcp
Environment="PATH=/home/preginald/DigitalSanctum/sanctum-mcp/venv/bin"
EnvironmentFile=/home/preginald/DigitalSanctum/sanctum-mcp/.env
ExecStartPre=/usr/bin/fuser -k -s 8100/tcp
ExecStart=/home/preginald/DigitalSanctum/sanctum-mcp/venv/bin/python server.py
ExecStopPost=/usr/bin/fuser -k -s 8100/tcp
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
