# ops/systemd/sanctum-worker.service
[Unit]
Description=Sanctum Signal Worker (Batch Processor)
After=network.target

[Service]
Type=oneshot
User=preginald
WorkingDirectory=/home/preginald/DigitalSanctum/sanctum-core
EnvironmentFile=/home/preginald/DigitalSanctum/sanctum-core/.env
ExecStart=/home/preginald/DigitalSanctum/sanctum-core/venv/bin/python3 -m app.worker

[Install]
WantedBy=multi-user.target